# server/systemd/start-servers.service
[Unit]
Description=Start services for home server

[Service]
Type=oneshot
ExecStart=/opt/rbin/start-servers
